# supervisord.conf
[supervisord] 
nodaemon=false

[program:nodejs] 
command=cd /app ; npm start