# supervisord.conf
[supervisord] 
nodaemon=true

[program:nodejs] 
directory=/app
command=npm start --production
autostart=true
autorestart=true

